# systemd/log-collector.timer
[Unit]
Description=Run Docker Log Collector Hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

